# unified-gateway.service
[Unit]
Description=Unified Gateway for Mission Control
After=network.target

[Service]
Type=simple
User=%I
WorkingDirectory=/home/%i/.openclaw/workspace/mission-control
ExecStart=/usr/bin/node unified-gateway.js
Environment="MEMORY_DIR=/home/%i/.openclaw/workspace/memory"
Environment="PATH=/usr/bin:/usr/local/bin"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
